# cron/monthly.borg.cron
#!/usr/bin/env bash

export BACKUP_SCHEDULED_JOB=1

BACKUP_SCRIPT_HOME="$(/usr/bin/dirname "$(/usr/bin/realpath "${BASH_SOURCE[0]}")")"

"${BACKUP_SCRIPT_HOME}/borg-run.sh" monthly