# samples/cron/cyborg.conf.cron
#!/usr/bin/env bash

CYBORG_HOME='/opt/cyborg/bin'