# samples/cron/weekly.borg.cron
#!/usr/bin/env bash

set -euo pipefail

BACKUP_SCRIPT_HOME="$(/usr/bin/dirname "$(/usr/bin/realpath "${BASH_SOURCE[0]}")")"

source "${BACKUP_SCRIPT_HOME}/cyborg.conf.cron"
export CYBORG_HOME
export CYBORG_FREQUENCY=weekly

"${BACKUP_SCRIPT_HOME}/run-borg.sh"